# elt/crontab
0 */12 * * * /usr/local/bin/python3 /code/AERO.py >> /var/log/cron.log 2>&1
